# system/https-proxy.service
[Unit]
Description=HTTPS Proxy

[Service]
ExecStart=/usr/local/sbin/https-proxy
